# automation/safehouse-youtube-budget.service
# SafeHouse YouTube daily watch-time budget daemon (long-running loop).
# Paths are rendered absolute by ../youtube-budget/install.sh on install.
# Runs as root: it writes the root-owned state file (/var/lib/safehouse) and splices
# the live Windows hosts file, and uses the Docker socket to read Pi-hole's FTL DB.
[Unit]
Description=SafeHouse YouTube daily watch-time budget enforcement daemon
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__REPO__
ExecStart=__DAEMON_SH__
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
